# deploy/supervisord.conf
[program:%(supervisor_program_name)s]
command=%(rungunicorn_script)s
user=%(django_user)s
stdout_logfile=%(supervisor_stdout_logfile)s
redirect_stderr=%(supervisor_redirect_stderr)s
